# deploy/systemd/vaultbase-mirror.timer
[Unit]
Description=Daily Vaultbase metadata mirror

[Timer]
OnCalendar=*-*-* 03:15:00
RandomizedDelaySec=20m
Persistent=true
Unit=vaultbase-mirror.service

[Install]
WantedBy=timers.target
