# deploy/systemd/vaultbase-retention.timer
[Unit]
Description=Weekly Vaultbase retention and repository prune

[Timer]
OnCalendar=Sun *-*-* 05:00:00
RandomizedDelaySec=30m
Persistent=true
Unit=vaultbase-retention.service

[Install]
WantedBy=timers.target
